# rootfs/etc/systemd/system/getty@tty1.service.d/noclear.conf
# Don't clear out boot messages before starting tty1
[Service]
TTYVTDisallocate=no
