# rootfs/etc/systemd/system/fake-hwclock.service.d/override.conf
# Force fake-hwclock to run after /boot is mounted
[Unit]
After=local-fs.target
Wants=local-fs.target
Requires=local-fs.target
